# cron/crontab
# Source of truth for kyle's user crontab. Checked into dotfiles so future-me
# (or anyone grepping) can find what's scheduled on this machine.
#
# Workflow
#   Inspect:  crontab -l                              (or: cat this file)
#   Edit:     $EDITOR ~/Code/dotfiles/cron/crontab
#   Install:  crontab ~/Code/dotfiles/cron/crontab    (overwrites live crontab)
#   Verify:   diff <(crontab -l) ~/Code/dotfiles/cron/crontab   (must be empty)
#
# Do NOT use `crontab -e`. It edits the live copy in-place, which drifts from
# this file and defeats the whole point of checking it in.
#
# Environment
#   Cron runs with a minimal PATH and no shell rc files. Source envato.sh for
#   tool paths (ruby, gh, jira, gws) and JIRA_API_TOKEN. Use absolute paths for
#   binaries and scripts.
#
# Time zones
#   Cron uses system time. Host is NZST (UTC+12), so NZ 11:00 = UTC 23:00 the
#   previous day. Weekdays NZ (Mon-Fri) = UTC (Sun-Thu) = cron dow 0-4.
#
# macOS gotcha
#   First time cron runs a job touching files outside /tmp, macOS may block it.
#   Grant Full Disk Access to /usr/sbin/cron in System Settings > Privacy.
#
# Columns: min hour dom mon dow  command

# --- envato-worklog ---------------------------------------------------------
# Backfill any missing weekday worklog from the last 7 days.
# Runs 11:00 NZST, Mon-Fri. Skips weekends and days already logged.
# Errors: tail /tmp/worklog-cron.log
# TODO: remove the osascript notification once trust is established.
0 23 * * 0-4 /bin/zsh -c 'source /Users/kyle/.config/direnv/envrc.d/envato.sh && log=""; for i in $(seq 1 7); do d=$(date -v-${i}d +\%Y-\%m-\%d); dow=$(date -v-${i}d +\%u); if [ "$dow" -le 5 ] && [ ! -f "$HOME/Worklog/$d.md" ]; then /Users/kyle/.rbenv/shims/ruby /Users/kyle/Code/Envato-non-market/AI-Dev-Tooling/claude-code/plugins/envato-worklog/skills/worklog/scripts/worklog.rb --date "$d" && log="$log wrote $d;" || log="$log FAILED $d;"; fi; done; osascript -e "display notification \"${log:-no missing weekdays}\" with title \"worklog cron\""' 2>> /tmp/worklog-cron.log
